# overlay/etc/systemd/system/ckb-watchdog.timer
[Unit]
Description=CKB Watchdog — check every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30

[Install]
WantedBy=timers.target
