# overlay/etc/systemd/system/ckb-health-check.timer
[Unit]
Description=CKB Health Check Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
