# overlay/etc/systemd/system/ckb-health-check.service
[Unit]
Description=CKB Node Health Check
After=ckb.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ckb-health-check
