# overlay/etc/systemd/system/ckb-watchdog.service
[Unit]
Description=CKB Application-Level Watchdog
After=ckb.service
Requires=ckb.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ckb-watchdog
