# overlay/etc/systemd/system/fiber-bridge.service
[Unit]
Description=Fiber HTTP Bridge
After=fiber-node.service

[Service]
Type=simple
WorkingDirectory=/opt/fiber-bridge
ExecStart=/usr/bin/node bridge.js
Restart=on-failure
RestartSec=5
Environment=PORT=9090
Environment=FIBER_RPC=http://localhost:8227

[Install]
WantedBy=multi-user.target
